# accounting_server/whatsapp.service
[Unit]
Description=WhatsApp API Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/accounting_server
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
Environment="DEBUG=whatsapp-web.js:*"
Environment="CHROMIUM_PATH=/usr/bin/google-chrome"

# إنشاء المجلدات المطلوبة
ExecStartPre=/bin/mkdir -p /root/accounting_server/logs
ExecStartPre=/bin/chmod 755 /root/accounting_server/logs
ExecStartPre=/bin/touch /root/accounting_server/logs/whatsapp_api.log
ExecStartPre=/bin/touch /root/accounting_server/logs/whatsapp_api.error.log
ExecStartPre=/bin/chmod 644 /root/accounting_server/logs/whatsapp_api.log
ExecStartPre=/bin/chmod 644 /root/accounting_server/logs/whatsapp_api.error.log

ExecStart=/usr/bin/node /root/accounting_server/whatsapp_server.js
Restart=always
RestartSec=10

# إعدادات السجلات
StandardOutput=append:/root/accounting_server/logs/whatsapp_api.log
StandardError=append:/root/accounting_server/logs/whatsapp_api.error.log

[Install]
WantedBy=multi-user.target 